# rootfs/etc/service/php-fpm/run
#!/bin/sh

exec 2>&1
source /etc/envvars

/usr/local/sbin/php-fpm --nodaemonize